# tvm-backup/borg.service
[Unit]
Description = Automated backup of important files

[Service]
Type = oneshot
ExecStart = bash /usr/lib/tvm-backup/run-backup
